# supervisord.conf
[supervisord]
nodaemon=true

[program:phpfpm]
command=php-fpm

[program:cron]
command=/bin/bash -c "crontab /etc/cron.d/jobs && cron -f"
